# services/offlineimap/run
#!/bin/bash

exec nice offlineimap -u Noninteractive.Basic 2>&1
