# services/urxvtd/run
#!/bin/bash

exec fghack urxvtd -q 2>&1
