# services/cpanmini_webserver/run
#!/bin/bash

/home/jon/perl/install/bin/minicpan_webserver 2>&1
