# services/emacs/run
#!/bin/bash

cd ~;
fghack emacs --daemon 2>&1
