# services/emacs/run
#!/bin/bash

fghack emacs --daemon 2>&1 
