# services/emacs/run
#!/bin/bash

cd ~;
exec fghack emacs --daemon 2>&1
